# projects/Dropbear/DEBIAN/postinst
#!/bin/sh

launchctl load /Library/LaunchDaemons/dropbear.plist

exit 0
